# alpine-nginx-php/rootfs/etc/supervisor/conf.d/nginx.conf
[program:nginx]
command=nginx -g 'daemon off;'
autorestart=true
startretries=3
startsecs=3
