# alpine-nginx-php/rootfs/etc/supervisor/conf.d/php-fpm.conf
[program:php-fpm]
command=php-fpm7 -F
autorestart=true
startretries=3
startsecs=3
